# workersctl
#!/bin/bash
#
# task-queue	Start up Task Queue
#
# chkconfig: 2345 99 01
# description: Task Queue runs jobs in parallel.

set -e
set -u

HOST=$(hostname -s)
if [[ $HOST =~ ^dev ]]; then
	ENV=dev
else
	ENV=prod
fi

CONFIG_FILE=/content/$ENV/rstar/etc/task-queue.sysconfig

if [ -f "$CONFIG_FILE" ]; then
    . $CONFIG_FILE
fi

TQHOME=${TQHOME:-/usr/local/dlib/task-queue}

TQUSER=${TQUSER:-rstar}

PIDFILE=${PIDFILE:-/var/lib/task-queue/taskqueueserver.pid}

LOGDIR=${LOGDIR:-/content/$ENV/rstar/tmp/mdi/task-queue/logs}

MQHOST=${MQHOST:-localhost}

USE_CGROUP=${USE_CGROUP:-true}

CGROUP_NAME=${CGROUP_NAME:-cpu:cpulimited}

FOREGROUND=${FOREGROUND:-false}

RETVAL=0

start() {
	PIDS=$(pgrep -f 'ruby.*workers.rb' | paste -sd' ')
	if [ -n "$PIDS" ]; then
		echo "Task Queue already running with pids '$PIDS'" 1>&2
		exit 1
	fi
	HOST_LOGDIR="$LOGDIR/$HOST"
	LOGFILE=$HOST_LOGDIR/worker.$HOST.log
	if [ ! -e "$HOST_LOGDIR" ]; then
		mkdir -m 0700 $HOST_LOGDIR
		chown $TQUSER $HOST_LOGDIR
	fi
	CMD="$TQHOME/workers.rb -m $MQHOST -l $LOGFILE -p $PIDFILE"
	if [ "$FOREGROUND" = "true" ]; then
		CMD="$CMD -f"
	fi
	if [ "$USE_CGROUP" = "true" ]; then
		CMD="cgexec -g $CGROUP_NAME $CMD"
	fi
	if [ "$TQUSER" = "$(whoami)" ]; then
		$CMD
	else
		if [ -f $LOGFILE ]; then
			TQGROUP=$(id -gn $TQUSER)
			chown $TQUSER:$TQGROUP $LOGFILE
		fi
		su -ls /bin/bash $TQUSER -c "$CMD"
	fi
	RETVAL=$?
}

stop() {
	if [ -f $PIDFILE ]; then
		kill `cat $PIDFILE`
		RETVAL=$?
	fi
}

restart() {
	stop
	sleep 5
	start
}

reload() {
	if [ -f $PIDFILE ]; then
		kill -HUP `cat $PIDFILE`
		RETVAL=$?
	fi
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    restart)
        restart
        ;;
    reload)
        reload
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|reload}" >&2
        exit 3
        ;;
esac

exit $RETVAL
