# doc/task-queue.cron
# /etc/cron.d/task-queue

# Shutdown server if shutdown file exists
*/5 * * * * rstar /usr/local/dlib/task-queue/shutdown.sh > /dev/null 2>&1

