# doc/task-queue.cron
# Check if Jenkins has deployed a new version of task-queue
# every 5 minutes
*/5 * * * * rstar /usr/local/lib/task-queue/check-update.sh /dev/null 2>&1

