# doc/task-queue.cron
# /etc/cron.d/task-queue

PATH=/usr/sbin:/sbin:/usr/bin:/bin

# Perform task-queue service actions based on commands read from tq-action.txt
# We check the file every 5 minutes
*/5 * * * * root /usr/local/dlib/task-queue/tq-cmd.sh > /dev/null 2>&1

# Update task-queue rpm package every Sunday at 3am
0 3 * * Sun root LOG_LEVEL=DEBUG /usr/bin/python3 /usr/local/dlib/task-queue/update-rpm.py > /dev/null 2>&1
