# doc/task-queue.cron
# /etc/cron.d/task-queue

# Restart server if restart trigger file exists
*/5 * * * * root /usr/local/dlib/task-queue/restart.sh > /dev/null 2>&1

# Shutdown server if stop trigger file exists
*/5 * * * * root /usr/local/dlib/task-queue/stop.sh > /dev/null 2>&1
